# test/cdh/supervisord.conf
[supervisord]
nodaemon=true

[program:hbase]
command=hbase master start

[program:opentsdb]
command=/tsdb/start_opentsdb.sh
