# docker/supervisor-serf.conf
[program:serf]
command=/opt/sei-bin/serf-start.sh
numprocs=1
autostart=true
autorestart=true

[program:serf-join]
command=/opt/sei-bin/serf-join.sh
autorestart=false
