# docker/supervisor-bosun.conf
[program:bosun]
command=/opt/bosun/bosun -c /opt/bosun/docker.conf
