# docker/supervisor-sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D

[program:sshd-create-key]
command=/opt/sei-bin/create_ssh_key.sh
autorestart=false
