# docker/data/supervisord-opentsdb.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/data/supervisord.pid

[program:hbase]
command=/hbase/start_hbase.sh
priority=10
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log

[program:opentsdb]
command=/tsdb/start_opentsdb.sh
priority=15
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
