# docker/data/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/data/supervisord.pid

[program:bosun]
command=/bosun/bosun -c /data/bosun.toml
priority=20
redirect_stderr=true
stdout_logfile=/dev/fd/1 ; redirects the process's stdout to supervisor's so that we can see it in the shell
stdout_logfile_maxbytes=0

[program:tsdbrelay]
command=/tsdbrelay/tsdbrelay %(ENV_TSDBRELAY_OPTS)s
priority=100
redirect_stderr=true
stdout_logfile=/dev/fd/1 ; redirects the process's stdout to supervisor's so that we can see it in the shell
stdout_logfile_maxbytes=0

[program:scollector]
command=/scollector/scollector -conf %(ENV_DATA_DIR)s/scollector.toml
priority=200
redirect_stderr=true
stdout_logfile=/dev/fd/1 ; redirects the process's stdout to supervisor's so that we can see it in the shell
stdout_logfile_maxbytes=0
